# ssh/dropbear/enableSSH_deb_src/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable ssh-server.service
systemctl start ssh-server.service

exit 0
